# callcoach-api.service
[Unit]
Description=NLJ Platform FastAPI Application
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/mnt/c/Users/aeroz/Documents/GitHub/nlj_viewer/backend
EnvironmentFile=/mnt/c/Users/aeroz/Documents/GitHub/nlj_viewer/backend/.env
ExecStart=/root/.local/bin/uv run uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target